# iaas_env/detection_rules/sop4.xml
<group name="atomic_red_team,">

  <rule id="100004" level="16">
    <if_group>sysmon_event1</if_group>

    <field name="win.eventdata.originalFileName" type="pcre2">(?i)schtasks\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)schtasks.*(/create).*?(onlogon|onstart|onstartup)</field>

    <description>Atomic Red Team T1053.005 #1 detected: schtasks.exe used to create OnLogon or OnStartup scheduled task</description>

    <mitre>
      <id>T1053.005</id>
    </mitre>

    <group>atomic_red_team,</group>
    <options>no_full_log</options>
  </rule>

</group>
